# com.iptest.autoupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.iptest.autoupdate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/andrewwang/projects/iptest/run_and_push.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/andrewwang/projects/iptest</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/andrewwang/projects/iptest/auto_run.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/andrewwang/projects/iptest/auto_run.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
